# target/linux/tr6560/norflash/base-files/etc/uci-defaults/10_prplmesh
[ ! -e /etc/config/prplmesh ] && exit 0

. /lib/functions/system.sh

set_mesh()
{
    local device=$1

    case $device in
        *)
            ;;
    esac
}

main()
{
    local device=$(board_name)

    # generate mesh related config
    `set_mesh $device`
}


main

exit 0
