# target/linux/tr6560/base-files/etc/hotplug.d/hostapd/02_hostapd_nft_qos
#!/bin/sh

. /lib/functions.sh

do_nftqos_wlan_config() 
{    
	local destinationMac sourceMac sourceIp

    config_get destinationMac $1 destinationMac
    config_get sourceMac $1 sourceMac
    config_get sourceIp $1 sourceIp

    [ "$destinationMac" != "$MAC" ] && [ "$sourceMac" != "$MAC" ] && return 1

    if [ -n "$destinationMac" ]; then
        [ "$destinationMac" != "$MAC" ] && [ "$sourceMac" != "$MAC" ] && return 1

        eval "/etc/init.d/nft-qos start 2>/dev/null"

    elif [ -n "$sourceIp" ]; then
        eval "/etc/init.d/nft-qos start 2>/dev/null"
    else
        return 1
    fi
}


do_sta()
{
    logger -t hostapd_nft_qos -p notice "ACTION=$ACTION" "MAC=$MAC"
	
    local qos_section="priority"
	
	config_load nft-qos
	config_foreach do_nftqos_wlan_config "$qos_section"
}

case "$ACTION" in
    "sta")
        do_sta $MAC $STATUS $SSID $IFNAME
    ;;
esac
